# server-build/box_templates/environment/etc/cron.d/run_one_offs
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""
HOME=/root/setup
BOX_NAME=${box_name}
03,18,33,48 * * * * root mkdir -p /root/setup ; cd /root/setup && aws s3 cp s3://${SERVICE_ASSETS_BUCKET}/scripts/common_functions.sh common_functions.sh ; source /root/setup/common_functions.sh ; run_one_offs

